# build/package/debian/etc/init.d/pilot-config
#!/bin/sh
### BEGIN INIT INFO
# Provides: pilot-config
# Required-Start: udev mountkernfs $remote_fs
# Required-Stop:
# Default-Start: S 2 3 4 5
# Default-Stop:
# Short-Description: Initializes PLC variables and pilot-node if installed
# Description:
### END INIT INFO

. /lib/lsb/init-functions

case "$1" in
  start)
    log_daemon_msg "Checking if PLC variables are present"
    if [ -e /etc/pilot/variables ]; then
      printf " Yes. Loading variables"
      cp /etc/pilot/variables /proc/pilot/plc/varconf/variables
      log_end_msg 0
    else
      log_end_msg 0
    fi
    ;;
  stop)
    ;;
  restart)
    ;;
  force-reload)
    ;;
  *)
    echo "Usage: $0 start" >&2
    exit 3
    ;;
esac
